# systemd/pyjob_server.conf
description "Server of the JobManager"
author "FAC"

setuid fac
setgid fac

start on net-device-up
stop  on runlevel [016]

task

exec /usr/local/bin/pyjob_server.py

console log
