# systemd/pyjob_server.conf
description "Server of the JobManager"
author "FAC"

setuid facs
setgid facs

start on net-device-up
stop  on runlevel [016]

task

exec pyjob_server.py

console log
